# Scripts/SwiftSnapshotsBot.service
[Unit]
Description=SwiftSnapshotsBot 

[Service]
Type=oneshot
ExecStart=/usr/local/bin/SwiftSnapshotsBot 

[Install]
WantedBy=multi-user.target
